# enable-ir-emitter.service
[Unit]
Description=enable ir emitter

[Service]
ExecStart=/usr/local/bin/enable-ir-emitter

[Install]
WantedBy=multi-user.target
