# boot_service/openrc/linux-enable-ir-emitter
#!/sbin/openrc-run

description="enable the infrared emitter"

start_pre() {
    modprobe uvcvideo
}

start() {
    linux-enable-ir-emitter run
}

depend() {
    need localmount
    after bootmisc
}
